# server/files/etc/cron.d/backup_app_db_to_s3
0                 6 * * *     app   /usr/local/ec2onrails/bin/backup_app_db.rb
05,15,25,35,45,55 * * * *     app   /usr/local/ec2onrails/bin/backup_app_db.rb --incremental
