# server/files/etc/cron.d/backup_app_db_to_s3
# Incremental backup every 5 minutes
*/5 * * * *  root  /usr/local/ec2onrails/bin/backup_app_db.rb --incremental

# Full backup every day at 05:01, reset the binary logs.
# First kill any incremental backup that happens to be in progress
1 5 * * *  root  killall -q -u root backup_app_db.rb ; /usr/local/ec2onrails/bin/backup_app_db.rb --reset
