# server/files/etc/cron.hourly/app
#!/bin/sh

test -e /mnt/app/current/script/hourly || exit 0
cd /mnt/app/current
/usr/local/ec2onrails/bin/rails_env sudo -u app script/hourly
