# server/files/etc/cron.weekly/app
#!/bin/sh

test -e /mnt/app/current/script/weekly || exit 0
cd /mnt/app/current
sudo -u app /usr/local/ec2onrails/bin/rails_env script/weekly
